# supervisor/tatris.ini
[program:tatris]
command=/home/admin/bin/start_tatris.sh
process_name=%(program_name)s
autostart=true

startretries=3
autorestart=unexpected
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
user=admin
redirect_stderr=true
stdout_logfile=/home/admin/stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=0
stdout_events_enabled=true
environment=USER="admin",HOME="/home/admin"
directory=/home/admin